# application/device/smartindustry/crontab
@reboot sleep 60 && /home/pi/go/src/github.com/smartfog/fogflow/application/device/smartindustry/start.sh & 

